# system_files/usr/lib/systemd/system/qualys-cloud-agent.service
[Unit]
Description=Qualys Cloud Agent
Documentation=https://www.qualys.com/
After=network-online.target
Wants=network-online.target
ConditionPathExists=/usr/libexec/qualys/cloud-agent/bin/qualys-cloud-agent

# Restart limits to prevent rapid restart loops
StartLimitIntervalSec=600
StartLimitBurst=3

[Service]
Type=simple
EnvironmentFile=-/etc/environment
EnvironmentFile=-/etc/default/qualys-cloud-agent
EnvironmentFile=-/etc/sysconfig/qualys-cloud-agent

# First-boot activation (idempotent - skips if already activated)
# Use '-' prefix to allow service to continue even if activation fails
# Timeout after 5 minutes to prevent hanging
ExecStartPre=-/bin/bash -c 'timeout 300 /usr/libexec/qualys/cloud-agent/bin/qualys-first-boot-activation.sh || true'

# Start the Qualys Cloud Agent
# Config.db writes are redirected to /var via symlink created by tmpfiles.d
ExecStart=/usr/libexec/qualys/cloud-agent/bin/qualys-cloud-agent

# Restart configuration
# Only restart on actual failure, not on clean exit
Restart=on-failure
RestartSec=60s
TimeoutStopSec=90s
TimeoutStartSec=120s

# Working directory
WorkingDirectory=/usr/libexec/qualys/cloud-agent

# Environment
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=HOME=/root

# Security settings (relaxed for Qualys agent requirements)
NoNewPrivileges=false
PrivateTmp=false
ProtectSystem=false
ProtectHome=false

# Allow access to necessary paths
ReadWritePaths=/var/log/qualys
ReadWritePaths=/var/lib/qualys
ReadWritePaths=/var/run/qualys
ReadWritePaths=/etc/qualys

[Install]
WantedBy=multi-user.target

